# ansible/template/supervisor.ini
[program:{{ app_name }}]
environment=OVH_APPLICATION_KEY="{{acme_api.OVH_APPLICATION_KEY}}",OVH_APPLICATION_SECRET="{{acme_api.OVH_APPLICATION_SECRET}}",OVH_CONSUMER_KEY="{{acme_api.OVH_CONSUMER_KEY}}",OVH_ENDPOINT="{{acme_api.OVH_ENDPOINT}}"
command=/usr/local/bin/traefik --configFile=config.toml
directory=/home/{{ app_name }}/
stdout_logfile=/home/{{ app_name }}/out.log
stderr_logfile=/home/{{ app_name }}/err.log
user=root
autostart=true
autorestart=true