# ansible/template/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/traefik --configFile=config.toml 
directory=/home/{{ app_name }}/
stdout_logfile=/home/{{ app_name }}/out.log
stderr_logfile=/home/{{ app_name }}/err.log
user=root
autostart=true
autorestart=true